# systemd/arkfile@.service
[Unit]
Description=Arkfile Service (%i instance)
After=network.target rqlite.target
Wants=caddy.service
Requires=rqlite.target

[Service]
Type=simple
User=ark%i
Group=arkfile
WorkingDirectory=/opt/arkfile/releases/current
EnvironmentFile=/opt/arkfile/etc/%i/secrets.env
ExecStart=/opt/arkfile/bin/arkfile

# Security measures
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Restart configuration
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
