# systemd/arkfile.service
[Unit]
Description=Arkfile Secure File Sharing Service
Documentation=https://github.com/84adam/arkfile
After=network.target
Wants=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
User=arkfile
Group=arkfile
WorkingDirectory=/opt/arkfile
Restart=always
RestartSec=5s
ExecStart=/opt/arkfile/bin/arkfile

# Configuration
EnvironmentFile=-/opt/arkfile/etc/secrets.env

# Library path for libopaque and liboprf shared libraries
Environment="LD_LIBRARY_PATH=/opt/arkfile/releases/current/vendor/stef/libopaque/src:/opt/arkfile/releases/current/vendor/stef/liboprf/src:/opt/arkfile/releases/current/vendor/stef/liboprf/src/noise_xk"

# Systemd credentials for sensitive keys
LoadCredential=opaque_server_key:/opt/arkfile/etc/keys/opaque/server_private.key
LoadCredential=jwt_signing_key:/opt/arkfile/etc/keys/jwt/current/signing.key

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/arkfile/var /opt/arkfile/etc/keys
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
