# systemd/arkfile.service
[Unit]
Description=Arkfile Secure File Sharing Service
Documentation=https://github.com/84adam/arkfile
After=network.target
Wants=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=arkfile
Group=arkfile
WorkingDirectory=/opt/arkfile/releases/current
ExecStart=/opt/arkfile/bin/arkfile

# Configuration
EnvironmentFile=-/opt/arkfile/etc/secrets.env

# Systemd credentials for sensitive keys
LoadCredential=opaque_server_key:/opt/arkfile/etc/keys/opaque/server_private.key
LoadCredential=jwt_signing_key:/opt/arkfile/etc/keys/jwt/current/signing.key

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Allow reads to configuration and writes to application directories
ReadOnlyPaths=/opt/arkfile/etc
ReadWritePaths=/opt/arkfile/var

# Restart policies
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
