# systemd/arkfile-test.service
[Unit]
Description=Arkfile Service (Test Instance)
After=network.target
Wants=caddy.service

[Service]
Type=simple
User=app
Group=app
WorkingDirectory=/opt/arkfile-test
ExecStart=/opt/arkfile-test/arkfile
Environment=HOST=test.arkfile.net
Environment=TEST_PORT=8081
Environment=BACKBLAZE_ENDPOINT=your_endpoint
Environment=BACKBLAZE_KEY_ID=your_key_id
Environment=BACKBLAZE_APPLICATION_KEY=your_app_key
Environment=BACKBLAZE_BUCKET_NAME=your_test_bucket
Environment=JWT_SECRET=your_test_jwt_secret

# Security measures
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes

# Restart configuration
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
