# systemd/minio@.service
[Unit]
Description=MinIO Server %i
Documentation=https://docs.min.io
After=network.target
PartOf=minio.target

[Service]
Type=simple
User=arkfile
Group=arkfile
WorkingDirectory=/opt/arkfile/var/lib/minio
EnvironmentFile=/opt/arkfile/etc/secrets.env

# Environment variables for MinIO server
Environment=MINIO_ROOT_USER=${MINIO_ROOT_USER}
Environment=MINIO_ROOT_PASSWORD=${MINIO_ROOT_PASSWORD}

# Server setup - simplified single-node mode
ExecStart=/usr/local/bin/minio server \
    --address :9000 \
    --console-address :9001 \
    ${LOCAL_STORAGE_PATH}

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
