# systemd/minio.service
[Unit]
Description=MinIO Server
Documentation=https://docs.min.io
After=network.target

[Service]
Type=simple
User=arkfile
Group=arkfile
WorkingDirectory=/opt/arkfile/var/lib/minio
EnvironmentFile=/opt/arkfile/etc/secrets.env

# Environment variables for MinIO server
Environment=MINIO_ROOT_USER=${MINIO_ROOT_USER}
Environment=MINIO_ROOT_PASSWORD=${MINIO_ROOT_PASSWORD}
Environment=MINIO_SSE_AUTO_ENCRYPTION=${MINIO_SSE_AUTO_ENCRYPTION}

# Server setup - single-node mode
ExecStart=/usr/local/bin/minio server \
    --address :9000 \
    --console-address :9001 \
    /opt/arkfile/var/lib/minio/data

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
