# systemd/rqlite@.service
[Unit]
Description=rqlite distributed database %i
Documentation=https://github.com/rqlite/rqlite
After=network.target
PartOf=rqlite.target

[Service]
Type=simple
User=arkfile
Group=arkfile
EnvironmentFile=/opt/arkfile/etc/secrets.env
WorkingDirectory=/opt/arkfile/var/lib/database
ExecStart=/usr/local/bin/rqlited \
    -http-addr :4001 \
    -raft-addr :4002 \
    -node-id demo \
    -auth /opt/arkfile/etc/rqlite-auth.json \
    data
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
